# govee-gateway.service
[Unit]
Description=Govee BLE Gateway
After=network.target bluetooth.target

[Service]
Type=simple
ExecStart=/opt/govee-gateway/venv/bin/python /opt/govee-gateway/gateway.py
WorkingDirectory=/opt/govee-gateway
Restart=always
RestartSec=10

# Load environment variables from file
EnvironmentFile=-/opt/govee-gateway/.env

# Bluetooth access
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
